# supervisord.ini
[supervisord]
nodaemon=true

[program:abrtd]
command=/usr/sbin/abrtd -d

[program:abrt-dump-journal-oops]
command=/usr/bin/abrt-dump-journal-oops -J /host/var/log/journal -fxtD

[program:abrt-dbus]
command=/usr/sbin/abrt-dbus -t 2147483647
