# supervisord.ini
[supervisord]
nodaemon=false

[program:dbus-daemon]
priority=0
command=/usr/bin/dbus-daemon --system --nofork

[program:abrtd]
priority=100
command=/usr/sbin/abrtd -d

[program:abrt-dbus]
priority=200
command=/usr/sbin/abrt-dbus -t 2147483647

[program:abrt-dump-journal-oops]
command=/usr/bin/abrt-dump-journal-oops -J /host/var/log/journal -fxtD
startsecs=0

[program:abrt-check-vmcores]
command=/usr/sbin/abrt-harvest-vmcore --kdump /host/etc/kdump.conf
startsecs=0

[program:abrt-container-coredump]
command=/usr/local/bin/abrt-container-coredump --install
startsecs=0
