# supervisord.ini
[supervisord]
nodaemon=true

[program:dbus-daemon]
command=/usr/bin/dbus-daemon --system

[program:abrtd]
command=/usr/sbin/abrtd -d

[program:abrt-dump-journal-oops]
command=/usr/bin/abrt-dump-journal-oops -J /host/var/log/journal -fxtD

[program:abrt-check-vmcores]
command=/usr/sbin/abrt-harvest-vmcore -r /host

[program:abrt-dbus]
command=/usr/sbin/abrt-dbus -t 2147483647

[program:abrt-container-coredump]
command=/usr/local/bin/abrt-container-coredump --install
