# minimal-fedora/etc/supervisord.ini
[supervisord]
nodaemon=true

[program:dbus-daemon]
priority=0
command=/usr/bin/dbus-daemon --system --nofork

[program:abrtd]
priority=100
command=/usr/local/bin/after-file /run/dbus/messagebus.pid /usr/sbin/abrtd -d

[program:abrt-dbus]
priority=200
command=/usr/local/bin/after-file /run/dbus/messagebus.pid /usr/sbin/abrt-dbus -t 2147483647

[program:abrt-dump-journal-oops]
command=/usr/local/bin/after-file /var/run/abrt/abrtd.pid /usr/bin/abrt-dump-journal-oops -J /host/var/log/journal -fxtD
autorestart=false
startsecs=0

[program:abrt-check-vmcores]
command=/usr/local/bin/after-file /var/run/abrt/abrtd.pid /usr/sbin/abrt-harvest-vmcore --kdump /host/etc/kdump.conf
autorestart=false
startsecs=0

[program:abrt-container-coredump]
command=/usr/local/bin/after-file /var/run/abrt/abrtd.pid /usr/local/bin/abrt-container-coredump --install
autorestart=false
startsecs=0
